# pi/ruuvitag.service.example
# systemd unit for the Ruuvitag MQTT publisher.
#
# Install:
#   sudo cp ruuvitag.service.example /etc/systemd/system/ruuvitag.service
#   # adjust paths and IOT_ENDPOINT, then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now ruuvitag.service

[Unit]
Description=Ruuvitag MQTT publisher (BLE -> AWS IoT Core)
After=network-online.target bluetooth.service
Wants=network-online.target

[Service]
Type=simple
User=root
# BLE scanning needs CAP_NET_ADMIN + CAP_NET_RAW. Running as root is
# the simplest option on a single-purpose Pi; on a multi-user box give
# CAP_NET_RAW + CAP_NET_ADMIN to the python interpreter via setcap.

Environment=IOT_ENDPOINT=YOUR-ENDPOINT-ats.iot.eu-north-1.amazonaws.com
Environment=IOT_CLIENT_ID=ruuvitag-pi-home
Environment=IOT_CERT_PATH=/etc/ruuvitag/cert.pem
Environment=IOT_KEY_PATH=/etc/ruuvitag/private.key
Environment=IOT_CA_PATH=/etc/ruuvitag/AmazonRootCA1.pem
Environment=PUBLISH_INTERVAL_SEC=60
Environment=SCAN_DURATION_SEC=5

ExecStart=/opt/ruuvitag-venv/bin/python /opt/ruuvitag/ruuvitag_publish.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
